# usr/lib/systemd/user/dialogbank.service
[Unit]
Description=Dialogbank

[Service]
ExecStart=xterm -e "/usr/share/python/bin/dbank 2>&1 | tee  ~/dialogbank.log"
Type=simple
Restart=always