# etc/xdg/autostart/dialogbank.desktop
[Desktop Entry]
Name=Dialogbank
Exec=systemctl --user start dialogbank.service
Type=Application