# deb/dialogbank.service
[Unit]
Description=Dialogbank

[Service]
ExecStart=/usr/share/python/bin/dbank
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target
